# k380-function-keys-conf-git/70-k380.rules
SUBSYSTEM=="input", ATTR{name}=="Keyboard K380", RUN+="/usr/bin/k380_enable_function_keys.sh $sys$devpath"
